# Pi/PyBlaster/cmake/debian/prerm
#!/bin/sh

update-rc.d pyblaster remove

# rm -rf /var/lib/pyblaster
